# docker/cron/crontab
# Rotate logs weekly on Sunday at midnight
0 0 * * 0 /app/docker/cron/rotate_logs.sh >> /app/logs/cron.log 2>&1
